# deploy/polymarket-bot-p1-gabagool.service
[Unit]
Description=Polymarket Bot Profile 1 Gabagool (tmux)
After=network.target litestream.service
Requires=litestream.service

[Service]
Type=oneshot
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/polymarket-bot
EnvironmentFile=/opt/polymarket-bot/.env

ExecStart=/usr/bin/bash -lc '/usr/bin/tmux has-session -t bot-p1 2>/dev/null || /usr/bin/tmux new-session -d -s bot-p1 -e "HOME=/home/ec2-user" '\''/opt/polymarket-bot/.venv/bin/python main.py run --markets BTC --durations 15 --profile 1 --strategy gabagool_dual --tag exp-1-dual-gabagol --dashboard; exec bash'\'''
ExecStop=/usr/bin/bash -lc '/usr/bin/tmux kill-session -t bot-p1 || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
